# amdgpu-bind-resume.service
[Unit]
Description=AMDGPU bind recovery on resume
DefaultDependencies=no
Before=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/true
ExecStop=/usr/local/bin/amdgpu-bind.sh resume
TimeoutStopSec=45
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sleep.target
